# scripts/com.talk.whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Optional LaunchAgent: runs the local Whisper server automatically at login.
     Install (opt-in):
       cp scripts/com.talk.whisper.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.talk.whisper.plist
     Uninstall:
       launchctl unload ~/Library/LaunchAgents/com.talk.whisper.plist
       rm ~/Library/LaunchAgents/com.talk.whisper.plist
     NOTE: the model path below is absolute — edit it if your home dir differs. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.talk.whisper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/whisper-server</string>
        <string>-m</string>
        <string>/Users/zakarnahapetyan/.talk/models/ggml-large-v3-turbo.bin</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/zakarnahapetyan/.talk/whisper-server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zakarnahapetyan/.talk/whisper-server.log</string>
</dict>
</plist>
